# ansible/templates/juicefs-meta-dump.service.j2
[Unit]
Description=Daily JuiceFS metadata dump for restic backup
After=juicefs-mount.service
Requires=juicefs-mount.service

[Service]
Type=oneshot
User=host
Group=host
# ``juicefs dump`` walks the SQLite metadata and writes a JSON blob
# capturing every inode + chunk pointer.  Restoring is two steps:
#   1. ``juicefs format`` against the SAME bucket
#   2. ``juicefs load`` the dump
# The dump itself is small (KBs to a few MBs); restic dedup makes
# storing it daily cheap.  Lands inside persistent_data_dir/openhost/
# so the backup app's existing rules already cover it.
#
# The path is parameterised on ``juicefs_dump_dir`` so an operator
# who has overridden ``data_root_dir`` in their config gets the dump
# in the matching location; the default matches the standard
# ansible-deployed path.
ExecStart=/bin/bash -c '\
    set -eu; \
    /usr/local/bin/juicefs dump \
        "sqlite3:///var/lib/juicefs/meta.db" \
        {{ juicefs_dump_dir | default("/home/host/.openhost/local_compute_space/persistent_data/openhost") }}/juicefs-metadata-dump.json.tmp; \
    mv {{ juicefs_dump_dir | default("/home/host/.openhost/local_compute_space/persistent_data/openhost") }}/juicefs-metadata-dump.json.tmp \
       {{ juicefs_dump_dir | default("/home/host/.openhost/local_compute_space/persistent_data/openhost") }}/juicefs-metadata-dump.json'
