# ansible/templates/juicefs-mount.service.j2
[Unit]
Description=JuiceFS mount for OpenHost app_archive tier
# openhost.service declares both After= AND Requires= on this unit
# when JuiceFS is enabled (see openhost.service.j2).  The Requires=
# is intentional and stronger than ordering: if this unit fails or
# stops, systemd refuses to start openhost.service, so apps can never
# silently write to the underlying empty mount-point and have those
# writes shadowed once the mount comes back.
After=network-online.target
Wants=network-online.target
Before=openhost.service
RequiresMountsFor=/var/lib/juicefs

[Service]
Type=forking
User=host
Group=host
# Secrets-only env file (0640 root:host).  ACCESS_KEY/SECRET_KEY are
# NOT passed on argv, so ``ps`` doesn't leak them.
EnvironmentFile=/etc/openhost/juicefs/s3.env
# JuiceFS reads ACCESS_KEY/SECRET_KEY from environment when --access-key
# isn't passed; ``mount`` itself has no creds-on-argv option that's
# safer than the env path, so this is the canonical pattern.
ExecStart=/usr/local/bin/juicefs mount \
    --background \
    --log /var/log/juicefs/mount.log \
    --no-usage-report \
    "sqlite3:///var/lib/juicefs/meta.db" \
    /data/app_archive_juicefs
# systemd doesn't run ExecStop through a shell, so a bare ``||`` would
# be passed as literal argv to ``juicefs umount``.  Wrap explicitly so
# the ``-l`` fallback runs when the clean unmount fails (e.g. because
# a container still has the FS busy and we're tearing the host down).
ExecStop=/bin/bash -c '/usr/local/bin/juicefs umount /data/app_archive_juicefs || /bin/umount -l /data/app_archive_juicefs'
# A failed mount must NOT silently fall through to local disk: if the
# mount drops, apps that wrote to /data/app_archive_juicefs/<app>/
# would start hitting the underlying empty mount point.  Restart loops
# would surface the issue but the safer behaviour is fail-loud.
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
