# ansible/templates/juicefs-meta-dump.timer.j2
[Unit]
Description=Daily JuiceFS metadata dump for restic backup

[Timer]
# Run every 24h with a randomised delay so restores can recover to
# a known-recent state even if the dump runs at a chaotic time.  An
# initial run 5 min after boot covers the case where the timer
# fires before the first periodic invocation.
OnBootSec=5min
OnUnitActiveSec=24h
RandomizedDelaySec=15min
Persistent=true

[Install]
WantedBy=timers.target
