# ansible/templates/openhost-juicefs.service.j2
[Unit]
Description=OpenHost JuiceFS Archive Mount
After=network-online.target
Wants=network-online.target
# Start before the main service so the FUSE mount is ready when
# compute_space boots and calls attach_on_startup.
Before=openhost.service

[Service]
Type=simple
User=host
# The mount command is backend-agnostic: JuiceFS reads the storage type +
# bucket from the meta DB (set by `format` / `config`).  Only the S3 backend
# needs credentials, which (when present) live as ACCESS_KEY / SECRET_KEY in
# the env file written by compute_space.  The unit stays installed-but-disabled
# until compute_space brings the archive backend online (at boot for the local
# file-backed volume, or on the dashboard S3-upgrade path); it then writes the
# env file and enables/starts this unit.
EnvironmentFile={{ juicefs_env_file }}
# systemd does not shell-expand $VAR in ExecStart; we use /bin/sh
# so the env vars from EnvironmentFile are expanded correctly.
ExecStart=/bin/sh -c '"$JUICEFS_BINARY" --no-agent mount --no-usage-report -o allow_other "$JUICEFS_META_DSN" "$JUICEFS_MOUNT_DIR"'
ExecStop=/bin/sh -c '"$JUICEFS_BINARY" umount "$JUICEFS_MOUNT_DIR"'
Restart=always
RestartSec=5
# JuiceFS logs to ~/.juicefs/juicefs.log; keep stdout/stderr quiet
# so journald doesn't duplicate everything.
StandardOutput=null
StandardError=journal

[Install]
WantedBy=multi-user.target
